# linux_config/wifi_check.timer
[Unit]
Description=Exécuter wifi_check.service toutes les minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target